# systemd/atlassian-ip-range-fw-cfg.timer
[Timer]
OnCalendar=weekly

[Install]
WantedBy=timer.target
