# systemd/atlassian-ip-range-fw-cfg.service
[Unit]
Description=Update Atlassian IP Range fw cfg

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/atlassian-ip-range-fw-cfg
